# tools/positioner-configs/dataproxy.service
[Unit]
Description=Dataproxy

[Service]
Type=simple
ExecStart=/usr/local/bin/dataproxy

[Install]
WantedBy=multi-user.target