# crow/system/crowker.service
[Unit]
Description=Crow brocker

[Service]
Type=simple
ExecStart=/usr/local/bin/crowker

[Install]
WantedBy=multi-user.target